# root/etc/services.d/squid/run
#!/usr/bin/with-contenv bash
/usr/sbin/squid -NYCd 5 -f /etc/squid/squid.conf
