# root/etc/s6-overlay/s6-rc.d/p10-squid-setup/up
/etc/s6-overlay/s6-rc.d/p10-squid-setup/run
