# root/etc/s6-overlay/s6-rc.d/svc-squid/run
#!/usr/bin/with-contenv bash
set -e

CONFDIR="${CONFDIR:-/etc/squid}";
SQUID_ARGS="${SQUID_ARGS:- -NYCd 5}";

exec \
    /usr/sbin/squid ${SQUID_ARGS} -f ${CONFDIR}/squid.conf
